# s6/btsync/run
#!/bin/sh

exec s6-applyuidgid -u $BTSYNC_UID -g $BTSYNC_GID \
    ./executable
