# s6/.s6-svscan/finish
#!/bin/sh

echo "s6 finished all services"
echo "doing last check..."

kill $( pidof rpc.mountd )
kill $( pidof rpcbind )
kill $( pidof rpc.statd )
kill $( pidof btsync )

echo "Done."
